# niri/udev/rules.d/61-dell-g5-5587-keyboard-remap.rules
ACTION=="add|change", SUBSYSTEM=="input", KERNEL=="event*", DRIVERS=="atkbd", ATTRS{name}=="AT Translated Set 2 keyboard", RUN+="/home/jake/.dotfiles/niri/remap-built-in-keyboard.sh %E{DEVNAME}"
